# deploy/debian_scripts/postrm
#!/bin/bash
set -e
case "$1" in
    remove)
        if command -v systemctl >/dev/null; then
            echo "Disabling kernelgatekeeper system service..."
            systemctl disable kernelgatekeeper.service || true
        fi
        # Remove profile.d script
        echo "Removing profile.d script..."
        rm -f /etc/profile.d/99-kernelgatekeeper-client-enabler.sh
        if command -v systemctl >/dev/null; then systemctl daemon-reload || true; fi
        ;;
    purge)
        echo "Purging kernelgatekeeper configuration, logs, and profile script..."
        rm -f /etc/kernelgatekeeper/config.yaml /etc/kernelgatekeeper/config.yaml.example
        rm -f /var/log/kernelgatekeeper.log
        rmdir /etc/kernelgatekeeper 2>/dev/null || true # Remove dir if empty
        rm -f /etc/profile.d/99-kernelgatekeeper-client-enabler.sh
        if command -v systemctl >/dev/null; then
            systemctl disable kernelgatekeeper.service || true
            systemctl daemon-reload || true
        fi
        ;;
    upgrade|abort*|disappear)
        ;;
    *)
        echo "postrm called with unknown argument '$1'" >&2
        exit 1
        ;;
esac
exit 0